# backend/config/apppanda-backend.service
[Unit]
Description=apppanda backend — HTTP/SSE bridge para la app Android AppPanda
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/apppanda-backend
# INSTALL.sh reemplaza __TARGET_USER__ y __TARGET_GROUP__ con el usuario actual.
User=__TARGET_USER__
Group=__TARGET_GROUP__
Restart=on-failure
RestartSec=5

# Estado: SQLite WAL + audit log. systemd crea los dirs.
StateDirectory=apppanda-backend
LogsDirectory=apppanda-backend
StateDirectoryMode=0700
LogsDirectoryMode=0750

# Hardening
# ProtectHome=no y RestrictNamespaces=no para que el subproceso `claude`
# spawneado por claude_runner pueda escribir en ~/.claude/ (session-env,
# cache, projects) y crear su propio mount namespace para su sandbox
# interno. Antes estaba ProtectHome=read-only y los hijos heredaban el
# bind-mount RO de /home, devolviendo EROFS al intentar inicializar.
# PrivateTmp también se afloja porque claude usa /tmp para artefactos
# intermedios que a veces tiene que compartir con apps GUI lanzadas.
ProtectSystem=strict
ProtectHome=no
ReadOnlyPaths=/etc/apppanda-backend
ReadWritePaths=/var/lib/apppanda-backend /var/log/apppanda-backend
PrivateTmp=no
NoNewPrivileges=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictNamespaces=no
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes

[Install]
WantedBy=multi-user.target
